# safety/supervisor/rpi/wu-kill-button.service
[Unit]
Description=Wu physical kill button (GPIO → killswitch)
Documentation=file:/home/kris/Development/Wu_v3/safety/supervisor.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Must run as a user that's in the `gpio` group (default `pi` works on
# Raspberry Pi OS, or add the chosen user to `gpio` manually).
User=pi
Group=pi
ExecStart=/usr/bin/python3 /opt/wu/wu_kill_button.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
